# dragongate.supervisor.conf
[program:dragongate]
directory=/root/dragongate
command=sudo -E -u root /usr/bin/nice -n 5 ./dragongate
redirect_stderr=true
