# dragongate.supervisor.conf
[program:dragongate]
directory=/app/dragongate
command=sudo -E -u root /usr/bin/nice -n 5 /app/go/src/github.com/dioptre/dragongate/dragongate
redirect_stderr=true
