# tracker.supervisor.conf
[program:tracker]
directory=/root/tracker
command=sudo -E -u root /usr/bin/nice -n 5 ./tracker
redirect_stderr=true
